# 30-unrealircd.conf

[program:unrealircd]
user=pissnet
command=/home/pissnet/unrealircd/bin/unrealircd -F
stopwaitsecs=30
