# 50-sshd.conf

[supervisord]
user=root

[program:sshd]
stopwaitsecs=30
command=/usr/sbin/sshd -eD
stderr_logfile=/dev/console
stderr_logfile_maxbytes=0
